# crontab
*/10 * * * * /app/venv/bin/python /app/monitor.py >> /var/log/cron.log 2>&1
