# crontab
*/10 * * * * root /usr/bin/python3 /app/monitor.py >> /var/log/cron.log 2>&1
